# docs/edutrack.service
[Unit]
Description=EduTrack Learning Management System
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/dotnet /opt/edutrack/EduTrack.WebApp.dll
WorkingDirectory=/opt/edutrack
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=edutrack
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

[Install]
WantedBy=multi-user.target
